# systemd/install.path
[Unit]
Description=Trigger for install.service

[Path]
PathExists=/home/install/ignition
Unit=install.service

[Install]
WantedBy=multi-user.target
